# anki_helper/anki_helper.timer
[Unit]
Description=Timer to run ank_helper daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
